# ActionBooster1/parallel-api.service
[Unit]
Description=ActionBooster Parallel API
After=network.target

[Service]
Type=simple
User=macbookm1
WorkingDirectory=/Users/macbookm1/Documents/ActionBooster1
ExecStart=/Users/macbookm1/Documents/ActionBooster1/ai_env/bin/python /Users/macbookm1/Documents/ActionBooster1/parallel_api.py
Restart=always
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
